# patchpilot.service
[Unit]
Description=PatchPilot Server
After=network.target

[Service]
Type=simple
User=patchpilot
Group=patchpilot
LogsDirectory=patchpilot
LogsDirectoryMode=0750
WorkingDirectory=/opt/patchpilot/server
EnvironmentFile=-/run/patchpilot-bootstrap.env
EnvironmentFile=/opt/patchpilot/.env
Environment=PATCHPILOT_DB_PATH=/opt/patchpilot/server/patchpilot.db
ExecStart=/opt/patchpilot/server/start.sh
Restart=always
RestartSec=5
# Hardening
NoNewPrivileges=yes
PrivateTmp=yes
ProtectHome=yes
# Allow pinging hosts for ping-only agent monitoring
AmbientCapabilities=CAP_NET_RAW
CapabilityBoundingSet=CAP_NET_RAW
# Resource limits
MemoryMax=512M
TasksMax=200
CPUQuota=80%

[Install]
WantedBy=multi-user.target
