# patchpilot.service
[Unit]
Description=PatchPilot Server
After=network.target

[Service]
Type=simple
User=patchpilot
Group=patchpilot
WorkingDirectory=/opt/patchpilot/server
EnvironmentFile=/opt/patchpilot/.env
ExecStart=/opt/patchpilot/server/start.sh
Restart=always
RestartSec=5
# Hardening
NoNewPrivileges=yes
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
